# launchd/Library/LaunchAgents/launched.llama_cpp_fim.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>launched.llama_cpp_fim</string>
    <key>ProgramArguments</key>
    <array>
      <string>sh</string>
      <string>-c</string>
      <string>/opt/homebrew/bin/llama-server -hf sweepai/sweep-next-edit-1.5B --port 8989 -ngl 99 -fa on -c 8192 -b 2048 -ub 2048 --cache-type-k q8_0 --cache-type-v q8_0 --mlock --prio 2</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
      <key>SuccessfulExit</key>
      <false/>
    </dict>
  </dict>
</plist>
